# phpadm svc/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = nginx -g 'daemon off;'
startretries=10
autorestart=true

[program:php-fpm7]
command=php-fpm7 -F -R
startretries=10
autorestart=true