# nginx svc/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command = /usr/sbin/sshd -D
startretries=10
autorestart=true

[program:nginx]
command=nginx -g 'daemon off;'
startretries=10
autorestart=true

